# buildroot/external/rootfs_overlay/lib/systemd/system/asbuiltreport-ovf-init.service
[Unit]
Description=AsBuiltReport Manager — OVF/guestinfo Property Initialisation
Documentation=https://www.asbuiltreport.com
After=local-fs.target
Before=network.target systemd-networkd.service
ConditionPathExists=!/var/lib/asbuiltreport/.ovf-init-done

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=60
StandardOutput=journal+console
StandardError=journal+console
ExecStart=/usr/local/bin/abr-ovf-init.sh

[Install]
WantedBy=sysinit.target
